# files/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm7.0]
command=/bin/bash -c "php-fpm7.0 -y /etc/php-fpm.conf -c /etc/php/7.0/php.ini"
autostart=true
autorestart=true
startretries=3

# Redirect output so we can see it using "docker logs"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
